# System/Archive/launchd-plists-2026-04-21/com.ericmanchester.thevault.overnight.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ericmanchester.thevault.overnight</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>bash /Users/ericmanchester/theVault/System/Scripts/preflight.sh &amp;&amp; cd /Users/ericmanchester/theVault &amp;&amp; source .venv/bin/activate &amp;&amp; python3 System/Scripts/overnight_processor.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/ericmanchester/theVault/System/Logs/overnight.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ericmanchester/theVault/System/Logs/overnight.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/ericmanchester/theVault</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ericmanchester/.gem/bin:/opt/homebrew/opt/ruby/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ericmanchester</string>
        <key>GEM_HOME</key>
        <string>/Users/ericmanchester/.gem</string>
        <key>GEM_PATH</key>
        <string>/Users/ericmanchester/.gem</string>
        <key>THEVAULT_CALENDAR_BACKEND</key>
        <string>auto</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
